# esuw.ini
[program:ESUW]
command=node /home/orsarfat/bin/ESUW/build/ESUW.server.bundle.js
autostart=yes
autorestart=yes